# debian/postinst
#!/bin/sh

update-rc.d clover defaults || exit 1
invoke-rc.d clover start || exit 1
